# systemd/datalogger-tuya.service.template
[Unit]
Description=Datalogueur Tuya Temperature Capture Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
ExecStart=__APP_DIR__/.venv/bin/python tools/tuya_temperature_capture.py --watch
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
